# telegram_bot/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python /app/bot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log
